# homebrew-automation/LaunchAgents/com.user.brew-update.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.brew-update</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>HOMEBREW_SCRIPTS_DIR/brew-update.sh</string>
    </array>
    
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    
    <key>StandardOutPath</key>
    <string>USER_HOME_DIR/Library/Logs/brew-update.log</string>
    
    <key>StandardErrorPath</key>
    <string>USER_HOME_DIR/Library/Logs/brew-update-error.log</string>
    
    <key>RunAtLoad</key>
    <false/>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>